# profile.d/activemq.csh
#!/bin/csh
setenv ACTIVEMQ_HOME /usr/share/activemq
