# profile.d/activemq.sh
#!/bin/sh
export ACTIVEMQ_HOME=/usr/share/activemq
