# profile.d/axis2.csh
#!/bin/csh
setenv AXIS2_HOME /home/slok/MDISS/distribuidos/resources/axis2-1.5.6/
