# profile.d/axis2.sh
#!/bin/sh
export AXIS2_HOME=/home/slok/MDISS/distribuidos/resources/axis2-1.5.6/
