# profile.d/hadoop.sh
#!/bin/sh
export HADOOP_INSTALL=/usr/share/hadoop
export PATH="${PATH}:${HADOOP_INSTALL}/bin:${HADOOP_INSTALL}/sbin"
