# profile.d/hadoop.csh
#!/bin/csh
setenv HADOOP_INSTALL /usr/share/hadoop
setenv PATH ${PATH}:${HADOOP_INSTALL}/bin:${HADOOP_INSTALL}/sbin
