# profile.d/jboss.csh
#!/bin/csh
setenv JBOSS_HOME /usr/share/jboss7
#setenv PATH ${PATH}:${JBOSS_HOME}}/bin

