# profile.d/jboss.sh
#!/bin/sh
export JBOSS_HOME=/usr/share/jboss7
#export PATH=$PATH:$JBOSS_HOME/bin
