# profile.d/maven.sh
#!/bin/sh
export M2_HOME=/usr/share/maven
export PATH=$PATH:$M2_HOME/bin
