# profile.d/node.csh
#!/bin/csh
setenv NODE_PATH /usr/lib64/node_modules
