# profile.d/scala.sh
#!/bin/sh
export SCALA_HOME=/usr/lib64/scala
export MANPATH="${MANPATH}:${SCALA_HOME}/man"
export PATH="${PATH}:${SCALA_HOME}/bin"
