# profile.d/tomcat.sh
#!/bin/sh
export CATALINA_HOME=/usr/share/tomcat7
export TOMCAT_HOME=$CATALINA_HOME
