# profile.d/tomcat.csh
#!/bin/csh
setenv CATALINA_HOME /usr/share/tomcat7
setenv TOMCAT_HOME ${CATALINA_HOME}
